# cloud.ini
[program:cloud]
command = /home/mohsen/PycharmProjects/jwt/conf.sh
redirect_stderr = true
stderr_logfile= /var/log/supervisor/celery_log
stdout_logfile= /var/log/supervisor/celery_log
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8
stopasgroup=true
autostart=true
autorestart=true
stopsignal=KILL
killasgroup=true
